# files/etc/systemd/system/pymedia.target
[Unit]
Description=Start pymedia scripts

[Install]
WantedBy=multi-user.target

