# files/home/pwbt/.config/systemd/user/bt-presence.target
[Unit]
Description=Starts/stops service bound to BT headset presence (should be managed by monitor-bt.py)
After=pipewire.service
